# app-hosting-linux/config-files/gunicorn.socket
# Gunicorn Socket Configuration
# Place this file at: /etc/systemd/system/gunicorn.socket
# Adjust paths and settings according to your project

[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn.sock

[Install]
WantedBy=sockets.target

